# services/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=0
loglevel=warn
user=root   

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700  

[program:palworld]
command=<replaced by init.sh>
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
killasgroup=true
stopsignal=INT
user=steam
events=PROCESS_STATE
environment=HOME="/home/steam/server/",USER="steam"

[eventlistener:palworld_events]
user=steam
command=/home/steam/server/services/listener/palworld_state_handler.sh
events=PROCESS_STATE
autostart=true
autorestart=true